# Source Code/Python-Scripts/daemons-wired/wired-syncing.service
[Unit]
Description=WIRED Syncing
After = network.target
[Service]
ExecStart= sudo -u orangepi python3 /home/orangepi/Desktop/SHELF3/WIRED-syncing.py
WorkingDirectory=/home/orangepi/Desktop/SHELF3
Restart = always
[Install]
WantedBy=multi-user.target
